# distro/detector-server-1.0/DEBIAN/postinst
#! /bin/sh
systemctl enable detector-server
systemctl start detector-server
true
