# distro/detector-server-1.0/DEBIAN/prerm
#! /bin/sh
systemctl stop detector-service
systemctl disable detector-service
true
